# main/units/j4status@.socket.in
[Unit]
Description=j4status stream socket

[Socket]
Service=j4status@%i.service
ListenStream=%t/@PACKAGE_NAME@/%i

[Install]
WantedBy=sockets.target
